# pve-login-idle-shutdown.service
[Service]
ExecStart=/usr/local/bin/pve-login-idle-shutdown.sh
